# setup/static/etc/init.d/mypi-control
#!/sbin/openrc-run

start()
{
    if [ -f /mypi-control/setup ]
    then
        /mypi-control/start-mypi-control
    fi
}
