# setup/static/etc/init.d/mypi-setup
#!/sbin/openrc-run

start()
{
    if [ -f /mypi-setup/setup ]
    then
        /mypi-setup/setup
    fi
}
